# src/main/resources/r2cloud.crontab
#!/bin/sh

# r2cloud is executed under user pi, thus unable to change cron.daily.
if [ -f /home/pi/r2cloud/DO_NOT_UPDATE ]; then
    printf "auto update is disabled" | systemd-cat -t r2cloud
	exit 0
fi

MAILTO=""

printf "updating r2cloud" | systemd-cat -t r2cloud

systemd-cat -t r2cloud apt-get -y update
systemd-cat -t r2cloud apt-get -y install r2cloud r2cloud-ui