# src/main/resources/r2cloud.crontab
#!/bin/sh

# r2cloud is executed under user pi, thus unable to change cron.daily.
if [ -f /home/pi/r2cloud/DO_NOT_UPDATE ]; then
	exit 0
fi

sudo apt-get update
sudo apt-get install r2cloud r2cloud-ui