# src/main/deb/postinst
chown ${config.user}:${config.group} ${config.installDir}/lib/*.jar
chmod 640 ${config.installDir}/lib/*.jar

cp ${config.installDir}/etc/r2cloud.service /etc/systemd/system/r2cloud.service;

systemctl enable r2cloud.service
systemctl start r2cloud.service

cp ${config.installDir}/etc/r2cloud-prod.conf /etc/nginx/nginx.conf;

cp ${config.installDir}/etc/r2cloud.crontab /etc/cron.daily/r2cloud
chmod 755 /etc/cron.daily/r2cloud

if [ ! -f ${config.installDir}/ssl/domain-chain.crt ]; then
    echo "Certificate is not found at ${config.installDir}/ssl/domain-chain.crt. Generating new"
    mkdir -p ${config.installDir}/ssl
    openssl genrsa -passout pass:x -out ${config.installDir}/ssl/server.pass.key 2048
    openssl rsa -passin pass:x -in ${config.installDir}/ssl/server.pass.key -out ${config.installDir}/ssl/domain.key
    rm ${config.installDir}/ssl/server.pass.key
    openssl req -new -key ${config.installDir}/ssl/domain.key -out ${config.installDir}/ssl/domain.csr -subj "/C=RU/ST=Moscow/L=Moscow/O=r2cloud/OU=r2cloud/CN=localhost"
    openssl x509 -req -days 365 -in ${config.installDir}/ssl/domain.csr -signkey ${config.installDir}/ssl/domain.key -out ${config.installDir}/ssl/domain-chain.crt
    rm ${config.installDir}/ssl/domain.csr
fi

if [ -f /run/nginx.pid ]; then
	nginx -s reload
fi

cp ${config.installDir}/etc/sudoers-nginx.conf /etc/sudoers.d/nginx
chmod 440 /etc/sudoers.d/nginx

usermod -a -G systemd-journal pi