# src/main/deb/postrm
if systemctl is-active r2cloud; then
	systemctl stop r2cloud.service
fi
if systemctl is-enabled r2cloud; then
	systemctl disable r2cloud.service
fi

# restore old nginx configuration if exists
if [ -f /etc/nginx/nginx.conf.backup ]; then
	cp /etc/nginx/nginx.conf.backup /etc/nginx/nginx.conf
fi