# src/main/deb/postrm
if systemctl is-active r2cloud; then
	systemctl stop r2cloud.service
fi
if systemctl is-enabled r2cloud; then
	systemctl disable r2cloud.service
fi

if [ -f /run/nginx.pid ]; then
	nginx -s reload
fi