# src/main/deb/postrm
if [ -e /etc/nginx/sites-enabled/instaseller ]; then
	rm /etc/nginx/sites-enabled/instaseller
	nginx -s reload
fi

if systemctl is-active r2cloud; then
	systemctl stop r2cloud.service
fi
if systemctl is-enabled r2cloud; then
	systemctl disable r2cloud.service
fi

rm -f /etc/systemd/system/r2cloud.service;

