# src/main/deb/postinst
chmod 755 ${config.installDir}/lib/download.sh
${config.installDir}/lib/download.sh ${config.installDir}/lib/ /usr/share/java/r2cloud/
chown -R ${config.user}:${config.group} ${config.installDir} /usr/share/java/r2cloud/*.jar
chmod 640 ${config.installDir}/lib/*.jar /usr/share/java/r2cloud/*.jar

if id "${config.user}" &>/dev/null; then
    echo "User '${config.user}' already exists."
else
    useradd -m ${config.user}
    echo "User '${config.user}' created."
fi

if ! getent group ${config.group} > /dev/null; then
    groupadd ${config.group}
    echo "Group '${config.group}' created."
fi

if [ $(getent group dialout) ]; then
    usermod -a -G dialout ${config.user}
fi

if [ $(getent group plugdev) ]; then
    usermod -a -G plugdev ${config.user}
fi

if [ $(getent group systemd-journal) ]; then
    usermod -a -G systemd-journal ${config.user}
fi

if [ -d /var/lib/nginx/proxy ]; then
	chown -R www-data:www-data /var/lib/nginx/proxy
fi

# try unload modules. if module in use, then modprobe return non-zero status. ignore it
modprobe --remove rtl2832 || true
modprobe --remove dvb_usb_rtl28xxu || true
modprobe --remove e4000 || true

if [ -d /run/systemd/system ]; then
	systemctl enable r2cloud.service
	systemctl start r2cloud.service
fi

# restore backup nginx configuration
if [ -f /etc/nginx/nginx.conf.backup ]; then
	mv /etc/nginx/nginx.conf.backup /etc/nginx/nginx.conf
fi
if [ -f /etc/nginx/sites-enabled/default ]; then
	rm /etc/nginx/sites-enabled/default
fi

if [ -f /etc/cron.daily/r2cloud ]; then
	chmod 755 /etc/cron.daily/r2cloud
fi

echo "version=${config.version}" > ${config.installDir}/etc/version.properties

if [ ! -f ${config.installDir}/ssl/domain-chain.crt ]; then
    echo "Certificate is not found at ${config.installDir}/ssl/domain-chain.crt. Generating new"
    mkdir -p ${config.installDir}/ssl
    openssl genrsa -passout pass:x -out ${config.installDir}/ssl/server.pass.key 2048
    openssl rsa -passin pass:x -in ${config.installDir}/ssl/server.pass.key -out ${config.installDir}/ssl/domain.key
    rm ${config.installDir}/ssl/server.pass.key
    openssl req -new -key ${config.installDir}/ssl/domain.key -out ${config.installDir}/ssl/domain.csr -subj "/C=RU/ST=Moscow/L=Moscow/O=r2cloud/OU=r2cloud/CN=raspberrypi.local"
    openssl x509 -req -days 358000 -in ${config.installDir}/ssl/domain.csr -signkey ${config.installDir}/ssl/domain.key -out ${config.installDir}/ssl/domain-chain.crt
    rm ${config.installDir}/ssl/domain.csr
fi

if [ -f /run/nginx.pid ]; then
	nginx -s reload
fi
