# src/main/deb/postrm
rm -f /etc/modprobe.d/r2cloud-blacklist.conf
rm -f /etc/cron.daily/r2cloud
rm -f /etc/sudoers.d/nginx
rm -f /home/pi/.wxtoimglic

if systemctl is-active r2cloud; then
	systemctl stop r2cloud.service
fi
if systemctl is-enabled r2cloud; then
	systemctl disable r2cloud.service
fi

rm -f /etc/systemd/system/r2cloud.service;

