# src/main/deb/postinst
chown ${config.user}:${config.group} ${config.installDir}/lib/*.jar
chmod 640 ${config.installDir}/lib/*.jar

cp ${config.installDir}/etc/r2cloud.service /etc/systemd/system/r2cloud.service;

systemctl enable r2cloud.service
systemctl start r2cloud.service

cp ${config.installDir}/etc/r2cloud-prod.conf /etc/nginx/nginx.conf;

cp ${config.installDir}/etc/r2cloud.crontab /etc/cron.daily/r2cloud
chmod 755 /etc/cron.daily/r2cloud

if [ -f /run/nginx.pid ]; then
	nginx -s reload
fi

cp ${config.installDir}/etc/sudoers-nginx.conf /etc/sudoers.d/nginx
chmod 440 /etc/sudoers.d/nginx

usermod -a -G systemd-journal pi